# roles/dns_threat_intel_sync/templates/sync-threat-intel.timer.j2
[Unit]
Description=Sync DNS threat intelligence feeds
Wants=network-online.target
After=network-online.target

[Timer]
OnBootSec=5min
OnUnitActiveSec=6h
Persistent=true

[Install]
WantedBy=timers.target
