# roles/metrics_uptime_check/templates/uptime-monitor.service.j2
[Unit]
Description=DNS Uptime Monitor
After=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/uptime-monitor.sh
Restart=always
User=nobody

[Install]
WantedBy=multi-user.target
